# package/custom/user-settings/files/etc/uci-defaults/15_nginx
#!/bin/sh

uci -q get nginx || exit 0

uci -q batch <<-EOF >/dev/null
	delete nginx._redirect2ssl.return
	set nginx._redirect2ssl.include='conf.d/*.locations'
	commit nginx
EOF

exit 0
