# package/kernel/mac80211/rtl8192cd/files/rtl8192cd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=11

start() {
	lsmod | grep -q '^rtl8192cd' && [ -x "/usr/local/sbin/load_calibration" ] && /usr/local/sbin/load_calibration
}
